# dotfiles/autostart/polkit-mate-authentication-agent-1.desktop
[Desktop Entry]
Name=PolicyKit Authentication Agent
Exec=/usr/libexec/polkit-mate-authentication-agent-1
Type=Application
NoDisplay=true
Hidden=true
